# deploy/oracle-vm/privatools-auto-deploy.timer
[Unit]
Description=Poll GitHub for PrivaTools deploys

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=15s
Persistent=true
Unit=privatools-auto-deploy.service

[Install]
WantedBy=timers.target
